# linux/MFServices/mfds.service
[Unit]
Description=Micro Focus Directory Server

[Service]
Type=simple
ExecStart=/home/support/MFSupport/MFServices/mfServices.sh mfds start
ExecStop=/home/support/MFSupport/MFServices/mfServices.sh mfds stop
ExecReload=/home/support/MFSupport/MFServices/mfServices.sh mfds restart
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target